# InvDash_WeeklyRefresh.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Weekly inventory dashboard refresh — regenerates PDF, Build/Burn chart, PPTX, restarts dashboard. Runs every Wednesday at 07:00 when new BQ data is available.</Description>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-07-08T07:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByWeek>
        <WeeksInterval>1</WeeksInterval>
        <DaysOfWeek>
          <Wednesday/>
        </DaysOfWeek>
      </ScheduleByWeek>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <ExecutionTimeLimit>PT2H</ExecutionTimeLimit>
    <Enabled>true</Enabled>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>C:\Users\r0c0jug\AppData\Local\Programs\Python\Python314\python.exe</Command>
      <Arguments>C:\Users\r0c0jug\OneDrive - Walmart Inc\Project\Trade Automation - Inventory\run_refresh.py</Arguments>
      <WorkingDirectory>C:\Users\r0c0jug\OneDrive - Walmart Inc\Project\Trade Automation - Inventory</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
